# bot/inotlisted-bot.service
[Unit]
Description=iNotListed Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=inotbot
Group=inotbot
WorkingDirectory=/home/inotbot/iNotListed
EnvironmentFile=/etc/inotbot/bot.env
Environment=MPLCONFIGDIR=/tmp/matplotlib-inotbot
ExecStart=/home/inotbot/.venv/bin/python -m bot.main
Restart=on-failure
RestartSec=10s

# Hardening — no listening sockets, no privileged work.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/inotbot/iNotListed/reports
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
